# svn-server/etc/services.d/subversion/run
#!/command/with-contenv sh

# From https://github.com/smebberson/docker-alpine/tree/master/alpine-apache

exec /usr/bin/svnserve -d --foreground -r /data/repositories --listen-port 3690;
