# root/etc/s6-overlay/s6-rc.d/init-svnrepo/run
#!/command/with-contenv sh
# Create SVN repository root

# /data/respositories: persistent storage for SVN repo
if [ ! -d /data/repositories ]; then
  install -d -o abc -g abc /data/repositories
fi
