# root/etc/s6-overlay/s6-rc.d/svc-svnserve/run
#!/command/with-contenv sh

exec s6-setuidgid abc \
  /usr/bin/svnserve \
    --daemon --foreground \
    --root /data/repositories \
    --listen-port 3690;