# root/etc/s6-overlay/s6-rc.d/init-apache/run
#!/usr/bin/with-contenv bash

chown -R abc:abc /var/www
chown -R abc:abc /var/log/apache2
chown -R abc:abc /run/apache2